# tests/io.docker.socket
[Unit]
Description=Docker Manager Socket

[Socket]
ListenStream=/var/run/docker.sock

[Install]
WantedBy=sockets.target
