# tests/io.podman.socket
[Unit]
Description=Pod Manager Socket

[Socket]
ListenStream=/run/podman/io.podman

[Install]
WantedBy=sockets.target
